# airflow-rabbitmq-docker/airflow.ini
[program:airflow-webserver]
command=su - airflow -c "airflow webserver"
directory=/usr/local/airflow
autostart=true
autorestart=true
startretries=3
stderr_logfile=/usr/local/webserver.err.log
stdout_logfile=/usr/local/webserver.log
user=root

[program:airflow-sch]
command=su - airflow -c "airflow scheduler" 
directory=/usr/local/airflow
autostart=true
autorestart=true
startretries=3
stderr_logfile=/usr/local/airflow/scheduler.err.log
stdout_logfile=/usr/local/airflow/scheduler.log
user=root

[program:airflow-worker]
command=su - airflow -c "airflow worker" 
directory=/usr/local/airflow
autostart=true
autorestart=true
startretries=3
stderr_logfile=/usr/local/airflow/worker.err.log
stdout_logfile=/usr/local/airflow/worker.log
user=root

